# cups-fuse/etc/systemd/user/cups-fuse@.path
[Unit]
Description=Starts mounting printers to %f, if this directory exists
PartOf=cups-fuse.service

[Path]
PathExists=%f
Unit=cups-fuse@%i.service
